# systemd/promnesia-watcher.path
[Unit]

[Path]
PathModified=/home/gtrun/Dropbox/org-notes/GTD
Unit=promnesia-restarter.service

