# systemd/promnesia-watcher.path
[Unit]

[Path]
PathModified=~/org-notes/
Unit=promnesia-restarter.service

