# systemd/promnesia-restarter.service
[Unit]

[Service]

ExecStart=~/.config/promnesia-nix/systemd/promnesia-restarter.sh
Type=oneshot

