# systemd/promnesia-restarter.service
[Unit]

[Service]

ExecStart=/home/test/.config/promnesia-nix/systemd/promnesia-restarter.sh
Type=oneshot

