# systemd/promnesia-watcher.path
[Unit]

[Path]
PathModified=/home/test/org-notes/
Unit=promnesia-restarter.service

